# linux/etc/init.d/kiosk
#!/sbin/openrc-run

name="kiosk"
command="/usr/bin/python3"
command_args="/home/kiosk/harddver/main.py"
command_background=true
pidfile="/run/kiosk.pid"
command_user="kiosk:kiosk"

start_pre() {
    export DISPLAY=:0
}