# init.d/kiosk
#!/sbin/openrc-run

name="kiosk"
command="/usr/bin/chromium"
command_args="--no-sandbox --kiosk --no-first-run --disable-infobars --noerrdialogs https://xecut-me.github.io/harddver/"
command_background=true
pidfile="/run/kiosk.pid"

depend() {
 need startx
}

start_pre() {
    export DISPLAY=:0
}